# server/scripts/mysqld-upstart.conf

description "mysqld"

start on runlevel [23]
stop on shutdown

exec /usr/bin/mysqld_safe >> /var/log/upstart-mysqld.log 2>&1

respawn

